# Prod/startupScripts/monkeyApp.conf
start on startup
exec forever start /App/app.js
